# x/libva-v4l2-request/profile.d/vaapi.sh
#!/bin/bash

export LIBVA_DRIVER_NAME=v4l2_request
